# bot.conf
[program:Karta]
directory=/home/Karta/
command=python3 bot.py

autostart=true
autorestart=true
stderr_logfile=/home/Karta/logfile_err.log
stdout_logfile=/home/Karta/logfile.log
